# system/supervisord.conf
[unix_http_server]
file=/tmp/supervisor.sock                       ; path to your socket file

[supervisord]
logfile=/dev/stdout
logfile_maxbytes=0
loglevel=error                                  ; info, debug, warn, trace
pidfile=/var/run/supervisord.pid                ; pidfile location
nodaemon=true                                  ; run supervisord as a daemon
minfds=1024                                     ; number of startup file descriptors
minprocs=200                                    ; number of process descriptors
user=root                                       ; default user
childlogdir=/tmp/

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock         ; use a unix:// URL  for a unix socket

[program:nginx]
command=/usr/sbin/nginx
stderr_logfile = /dev/stdout
stdout_logfile = /dev/stderr
stdout_logfile_maxbytes = 0
stderr_logfile_maxbytes = 0

[program:uwsgi]
command=uwsgi --http :9090 --mount /=server:app --manage-script-name --threads=8
user=www-data
directory=/code
stderr_logfile = /dev/stdout
stdout_logfile = /dev/stderr
stdout_logfile_maxbytes = 0
stderr_logfile_maxbytes = 0

[program:uwsgi_reloader]
command=bash -c 'while true; do inotifywait -e modify,attrib,create,delete -r -q /api_src; rsync -a /api_src/ /code/; pkill -HUP uwsgi; done'
user=www-data
directory=/code
stderr_logfile = /dev/stdout
stdout_logfile = /dev/stderr
stdout_logfile_maxbytes = 0
stderr_logfile_maxbytes = 0

[program:cert_reloader]
command=bash -c 'while true; do inotifywait -e modify,attrib,create,delete -q /config; rsync -a /config/*.pem /code/; pkill -HUP nginx; done'
user=root
directory=/code
stderr_logfile = /dev/stdout
stdout_logfile = /dev/stderr
stdout_logfile_maxbytes = 0
stderr_logfile_maxbytes = 0

[program:extra_loader]
command=bash -c 'test -d /config/updates && while true; do inotifywait -e modify,attrib,create,delete -r -q /config/updates; rsync -a /config/updates/ /app/dist/; done'
user=www-data
directory=/tmp
stderr_logfile = /dev/stdout
stdout_logfile = /dev/stderr
stdout_logfile_maxbytes = 0
stderr_logfile_maxbytes = 0
